# system/rc.local
#!/bin/sh
# /etc/rc.local

powertop --auto-tune
